# config/com.marketintel.youtube.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the YouTube/Podcasts collector. Runs the local wrapper once
  a day from a residential IP (Actions IPs are bot-blocked by YouTube — see
  scripts/run_youtube_local.sh). Copy to ~/Library/LaunchAgents/ and load:

    cp config/com.marketintel.youtube.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.marketintel.youtube.plist

  ProgramArguments points at the collector inside the Dropbox working repo
  (.../market-intelligence/market-intelligence/scripts/run_youtube_local.sh).
  The old standalone clone (~/market-intel-youtube) is deprecated/removed — do
  not point here. Change the /Users/bacopa paths if your username or repo
  location differ. launchd cannot expand $HOME, so paths are absolute. If the
  Mac is asleep at the scheduled time, launchd runs the job once on wake.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marketintel.youtube</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bacopa/Dropbox (Personal)/Business/InvestTool/market-intelligence/market-intelligence/scripts/run_youtube_local.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>12</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-youtube.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-youtube.log</string>
</dict>
</plist>
