# config/com.marketintel.concepts-backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the Concepts historical backfill collector.
  SEPARATE from the Concepts DAILY stream (cloud GitHub Actions cron, 10:00 UTC).
  NO git operations anywhere. Writes only to concepts-history/ (outside the repo).

  Fires at 16:45 and 21:45 local — the owner's at-computer window, mirroring the
  YouTube backfill (16:00 / 21:00) but offset +45 min so the two heavy local jobs
  never start together:
    • Concepts daily cloud cron   10:00 UTC (= 12:00 local)
    • YouTube daily collector     12:00 local (launchd)
    • YouTube backfill            16:00 / 21:00 local (launchd)
    • Concepts backfill           16:45 / 21:45 local (this job)
  UNCAPPED: one firing drains the whole Stage-1+2 signal-tier backlog in a single
  multi-hour run (~30 different sites, load distributed). Resumable via file-dedup
  (+ _state/backfill_state.json) — the 21:45 firing / next day continues any tail.
  Per-SITE politeness (spaced same-domain fetches) and the hard-stop / 24h backoff
  on any per-site 429 / bot-block are retained.

  The wrapper (scripts/run_concepts_backfill_local.sh) manages:
    • PID lockfile (single instance)
    • 24h BACKOFF after any exit-10 (429 / bot-check / consecutive fails)
    • osascript notification on hard-stop / unexpected exit

  Logs:
    ~/Dropbox (Personal)/Business/InvestTool/market-intelligence/
      concepts-history/_runlog/launchd_YYYY-MM-DD.md   (dated session log)
    ~/Library/Logs/market-intel-concepts-backfill.log    (launchd stream)

  INSTALL:
    cp config/com.marketintel.concepts-backfill.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.marketintel.concepts-backfill.plist
  UNLOAD:
    launchctl unload ~/Library/LaunchAgents/com.marketintel.concepts-backfill.plist
  RUN NOW (test):
    launchctl start com.marketintel.concepts-backfill

  NOTE: launchd cannot expand $HOME — paths are absolute /Users/bacopa.
  Update if the username or clone location differs.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marketintel.concepts-backfill</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bacopa/Dropbox (Personal)/Business/InvestTool/market-intelligence/market-intelligence/scripts/run_concepts_backfill_local.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>16</integer>
      <key>Minute</key><integer>45</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>21</integer>
      <key>Minute</key><integer>45</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-concepts-backfill.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-concepts-backfill.log</string>
</dict>
</plist>
