# config/com.marketintel.concepts-residential.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the Concepts DAILY residential supplement — recovers
  cloud-IP-blocked sources (Chatham House, IEA Oil Market Report) that the
  10:00 UTC GitHub Actions cron can't reach, mirroring how
  com.marketintel.youtube.plist recovers YouTube from Actions' block. Writes
  into the tracked raw/concepts/ tree via scripts/collect_concepts.py --ids
  chatham_house,iea_oil_market_report --residential, then commits + pushes
  straight to main so the extraction routine picks the items up the normal
  way. Does NOT touch concepts-history/ or the history-backfill state.

  NOT INSTALLED BY DEFAULT — this file only exists in the repo so the owner
  can review and install it deliberately:

    cp config/com.marketintel.concepts-residential.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.marketintel.concepts-residential.plist

  Scheduled for 12:15 local — after the YouTube collector (12:00) and clear of
  the Concepts daily cloud cron (10:00 UTC = 12:00 local), the Concepts HEAVY
  local jobs (17:30/22:30), and the Concepts signal backfill (16:45/21:45).
  Runs against only 2 sources with a tight --days window, so it finishes in
  well under a minute — timing precision doesn't matter much here, this slot
  is just clear of every other local job.

  UNLOAD:
    launchctl unload ~/Library/LaunchAgents/com.marketintel.concepts-residential.plist

  RUN NOW (manual test, does not require installing the plist):
    bash scripts/run_concepts_residential_local.sh

  NOTE: launchd cannot expand $HOME — paths are absolute /Users/bacopa.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marketintel.concepts-residential</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bacopa/Dropbox (Personal)/Business/InvestTool/market-intelligence/market-intelligence/scripts/run_concepts_residential_local.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>12</integer>
    <key>Minute</key><integer>15</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-concepts-residential.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-concepts-residential.log</string>
</dict>
</plist>
