# config/com.marketintel.youtube-backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the YouTube historical backfill collector.
  SEPARATE from com.marketintel.youtube (daily 12:00). NO git operations anywhere.
  Fires at 16:00 and 21:00 local — staggered from the daily collector, same residential IP.

  The wrapper (scripts/run_youtube_backfill_local.sh) manages:
    • flock-guarded single instance (prevents stacked runs)
    • 24-h BACKOFF after any exit-10 (IP/bot signal)
    • osascript notifications on hard-stop or unexpected exit

  Logs:
    ~/Dropbox (Personal)/Business/InvestTool/market-intelligence/
      youtube-history/_runlog/launchd_YYYY-MM-DD.md   (dated session log, via tee)
    ~/Library/Logs/market-intel-youtube-backfill.log   (launchd stream capture)

  INSTALL:
    cp config/com.marketintel.youtube-backfill.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.marketintel.youtube-backfill.plist

  UNLOAD:
    launchctl unload ~/Library/LaunchAgents/com.marketintel.youtube-backfill.plist

  RUN IMMEDIATELY (test without waiting for 16:00):
    launchctl start com.marketintel.youtube-backfill

  NOTE: launchd cannot expand $HOME — all paths are absolute /Users/bacopa.
  Update if username or clone location differs.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marketintel.youtube-backfill</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bacopa/Dropbox (Personal)/Business/InvestTool/market-intelligence/market-intelligence/scripts/run_youtube_backfill_local.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>16</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>21</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-youtube-backfill.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-youtube-backfill.log</string>
</dict>
</plist>
