# config/com.marketintel.concepts-heavy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the Concepts HEAVY think-tank backfill.
  SEPARATE from the signal-tier backfill (com.marketintel.concepts-backfill) and the
  Concepts DAILY stream (cloud GitHub Actions). NO git operations anywhere. Writes
  only to concepts-history/ (outside the repo) via its own _state/heavy_state.json.

  Fires at 17:30 and 22:30 local — staggered from every other heavy local job:
    • YouTube daily collector     12:00 local (launchd)
    • Concepts daily cloud cron   10:00 UTC (= 12:00 local)
    • YouTube backfill            16:00 / 21:00 local (launchd)
    • Concepts signal backfill    16:45 / 21:45 local (launchd)
    • Concepts HEAVY              17:30 / 22:30 local (this job)
  Thousands of fetches: each firing drains a chunk of the 15 heavy sources through
  the structural genre+length filter (+render-retry), resumable via file-dedup +
  heavy_state.json done-set, with per-SITE quarantine (no global halt), then NO-OPs
  once complete.

  The wrapper (scripts/run_concepts_heavy_local.sh) manages the PID lockfile, the
  per-site quarantine surfacing, and osascript notifications.

  Logs:
    ~/Dropbox (Personal)/Business/InvestTool/market-intelligence/
      concepts-history/_runlog/launchd_heavy_YYYY-MM-DD.md   (dated session log)
    ~/Library/Logs/market-intel-concepts-heavy.log            (launchd stream)

  INSTALL:
    cp config/com.marketintel.concepts-heavy.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.marketintel.concepts-heavy.plist
  UNLOAD:
    launchctl unload ~/Library/LaunchAgents/com.marketintel.concepts-heavy.plist
  RUN NOW (first run, then it continues OS-managed):
    launchctl start com.marketintel.concepts-heavy

  NOTE: launchd cannot expand $HOME — paths are absolute /Users/bacopa.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marketintel.concepts-heavy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bacopa/Dropbox (Personal)/Business/InvestTool/market-intelligence/market-intelligence/scripts/run_concepts_heavy_local.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>17</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>22</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-concepts-heavy.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bacopa/Library/Logs/market-intel-concepts-heavy.log</string>
</dict>
</plist>
